# ansible/roles/zigbee2mqtt/templates/zigbee2mqtt/zigbee2mqtt.service.j2
[Unit]
Description=zigbee2mqtt
After=network.target

[Service]
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm start
WorkingDirectory=/opt/zigbee2mqtt
StandardOutput=null
# Or use StandardOutput=null if you don't want Zigbee2MQTT messages filling syslog, for more options see systemd.exec(5)
StandardError=null
Restart=always
RestartSec=10s
User={{ host_alt_user }}

[Install]
WantedBy=multi-user.target